# wireguard-ui/rootfs/etc/services.d/wireguard-ui/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2191
# ==============================================================================
# Home Assistant Community Add-on: Wireguard UI
# Runs the Wireguard UI
# ==============================================================================

declare clear_session_on_boot
declare dns
declare endpoint
declare port
declare single_user
declare user_loglevel
declare user_loglevel
declare subnet_ranges
declare default_interface_addresses
# WGUI_SERVER_INTERFACE_ADDRESSES="10.252.1.0/24,10.252.2.0/24" -e SUBNET_RANGES="Home:10.252.1.0/24; DMZ Zone: 10.252.2.0/24"

bashio::log.info 'Starting Wireguard UI...'

clear_session_on_boot=$(bashio::config 'clear_session_on_boot' || echo '')
dns=$(bashio::config 'default_client_dns' || echo '')
endpoint=$(bashio::config 'default_client_endpoint' || echo '')
port=$(bashio::addon.port '51820/udp' || echo '')
single_user=$(bashio::config 'single_user_mode' || echo '')
user_loglevel=$(bashio::config 'log_level' || echo '')
subnet_ranges=$(bashio::config 'subnet_ranges' || echo '')
default_interface_addresses=$(bashio::config 'default_interface_addresses' || echo '')

# settings directly passed to WireGuard UI without additional validation
EMAIL_FROM_ADDRESS=$(bashio::config 'email_from_address' || echo '')
EMAIL_FROM_NAME=$(bashio::config 'email_from_name' || echo '')
SENDGRID_API_KEY=$(bashio::config 'sendgrid_api_key' || echo '')
SMTP_HOSTNAME=$(bashio::config 'smtp_hostname' || echo '')
SMTP_PORT=$(bashio::config 'smtp_port' || echo '')
SMTP_USERNAME=$(bashio::config 'smtp_username' || echo '')
SMTP_PASSWORD=$(bashio::config 'smtp_password' || echo '')
SMTP_AUTH_TYPE=$(bashio::config 'smtp_auth_type' || echo '')
SMTP_ENCRYPTION=$(bashio::config 'smtp_encryption' || echo '')
SMTP_HELO=$(bashio::config 'smtp_helo' || echo '')
TELEGRAM_TOKEN=$(bashio::config 'telegram_token' || echo '')
TELEGRAM_ALLOW_CONF_REQUEST=$(bashio::config 'telegram_allow_conf_request' || echo '')
TELEGRAM_FLOOD_WAIT=$(bashio::config 'telegram_flood_wait' || echo '')

# wireguard-ui will SHA512 checksum the contents of the random password and use
# that as the password to encrypt sessions.
random_password() (
  set +eo pipefail
  local chars='-;.~,.<>[]{}!@#$%^&*()_+=`0-9a-zA-Z'
  if [ -n "${1:-}" ]; then
    chars="${1}"
  fi
  # tr will exit with 141 because its stdout is cut short by head
  tr -dc -- "${chars}" 2> /dev/null < /dev/urandom | \
    head -c128 || [ "$?" -eq 141 ]
  echo
)

# Run wireguard-ui
# wireguard-ui creates databases in `./db`
# /data is a persistent path in the container
cd /data
# missing state dir causes login loops
[ -d db ] || mkdir db
if [ "${clear_session_on_boot:-}" = true ]; then
  rm -f session-key
fi
if [ ! -f session-key ] || [ -z "$(<session-key)" ]; then
  random_password > session-key
  echo "Generated new session encryption key." >&2
fi
export SESSION_SECRET_FILE=session-key
# customizations from home assistant
WGUI_DEFAULT_CLIENT_ALLOWED_IPS="0.0.0.0/0,::/0"
if [ -z "${dns:-}" ]; then
  WGUI_DNS="$(ip route | awk '$1 == "default" { print $3 }')"
else
  WGUI_DNS="${dns}"
fi
WGUI_ENDPOINT_ADDRESS="${endpoint}:${port%/*}"
case "${user_loglevel:-}" in
  debug|info|warn|error|off)
    WGUI_LOG_LEVEL="${user_loglevel}"
    ;;
  *)
    WGUI_LOG_LEVEL=debug
    ;;
esac
SUBNET_RANGES="${subnet_ranges:-}"
WGUI_MTU=1420
WGUI_PERSISTENT_KEEPALIVE=0
WGUI_SERVER_INTERFACE_ADDRESSES="${default_interface_addresses}"
WGUI_SERVER_LISTEN_PORT="${port%/*}"
export SUBNET_RANGES WGUI_DEFAULT_CLIENT_ALLOWED_IPS WGUI_DNS
export WGUI_ENDPOINT_ADDRESS WGUI_LOG_LEVEL WGUI_MTU WGUI_PERSISTENT_KEEPALIVE
export WGUI_SERVER_INTERFACE_ADDRESSES WGUI_SERVER_LISTEN_PORT

# contact settings set from home assistant config
export EMAIL_FROM_ADDRESS EMAIL_FROM_NAME SENDGRID_API_KEY SMTP_AUTH_TYPE
export SMTP_ENCRYPTION SMTP_HELO SMTP_HOSTNAME SMTP_PASSWORD SMTP_PORT
export SMTP_USERNAME TELEGRAM_ALLOW_CONF_REQUEST TELEGRAM_FLOOD_WAIT
export TELEGRAM_TOKEN
if [ -f /applied_patches ]; then
  cat /applied_patches
fi
wgui_args=( -bind-address "127.0.0.1:8080" )
if [ "${single_user:-}" = true ]; then
  wgui_args+=( -disable-login )
fi
exec /wireguard-ui "${wgui_args[@]}"
